# part2/myshell.service
[Unit]
Description=My COSC330 Shell Service
After=network.target

[Service]
Type=simple
# Update this path to match where you deploy the shell binary on EC2
ExecStart=/home/ec2-user/part1/shell
WorkingDirectory=/home/ec2-user/part1
StandardInput=null
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
